# target/linux/rockchip/armv8/base-files/etc/init.d/pwm-fan
#!/bin/sh /etc/rc.common

START=21
start() {
    /usr/bin/pwm-fan.sh &
}
