# target/linux/rockchip/armv8/base-files/etc/init.d/resize-rootfs
#!/bin/sh /etc/rc.common

START=22
start() {
	rm -rf /etc/rc.d/S22resize-rootfs
	#service resize-rootfs disable
	/usr/bin/resize-rootfs.sh
}

