# target/linux/rockchip/armv8/base-files/etc/init.d/typec_usb2
#!/bin/sh /etc/rc.common

START=21
start() {
    /usr/bin/typec_usb2.sh
}
